# ansible/roles/botkit/templates/botkit.ini
[program:botkit]
autorestart = True
log_stderr = True
directory = /var/www/slackbot
user = {{ username }}
stopsignal = QUIT
log_stdout = True
environment=SLACK_TOKEN="{{ botkit_token }}",USER="{{ username }}",HOME="/home/{{ username }}"
command = /bin/node /var/www/slackbot/app.js
autostart = True
stderr_logfile = /var/log/supervisor/botkit.err.log
stdout_logfile = /var/log/supervisor/botkit.stdout.log
